# google-disk-expand.init
#!/bin/bash

### BEGIN INIT INFO
# Provides:          google-disk-expand
# Required-Start:    $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Resize the root filesystem if needed
### END INIT INFO

/usr/bin/google_disk_expand
